# 99-mouse-imwheel.rules
ACTION=="add|remove", SUBSYSTEM=="input", RUN+="/bin/su nicola -c '/home/ryan/dotfiles/imwheel-mouse-manager.sh'"

